# src/packaging/ui-slim/prerm
#!/bin/sh
# Project: Crankshaft
# This file is part of Crankshaft project.
# Copyright (C) 2025 OpenCarDev Team
#
#  Crankshaft is free software: you can redistribute it and/or modify
#  it under the terms of the GNU General Public License as published by
#  the Free Software Foundation; either version 3 of the License, or
#  (at your option) any later version.
#
#  Crankshaft is distributed in the hope that it will be useful,
#  but WITHOUT ANY WARRANTY; without even the implied warranty of
#  MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
#  GNU General Public License for more details.
#
#  You should have received a copy of the GNU General Public License
#  along with Crankshaft. If not, see <http://www.gnu.org/licenses/>.

set -e
set -x  # Debug: trace execution to find where prerm terminates
export DEBIAN_FRONTEND=noninteractive

case "$1" in
    remove|deconfigure)
        echo "[prerm ui-slim] phase=remove"
        # Stop and disable systemd service (if systemd is present)
        if command -v systemctl >/dev/null 2>&1 && [ -d /run/systemd/system ]; then
            echo "[prerm ui-slim] stopping service"
            systemctl stop crankshaft-ui-slim.service >/dev/null 2>&1 || true
            echo "[prerm ui-slim] stopping display setup service"
            systemctl stop crankshaft-ui-slim-display-setup.service >/dev/null 2>&1 || true
            echo "[prerm ui-slim] disabling service"
            systemctl disable crankshaft-ui-slim.service >/dev/null 2>&1 || true
            echo "[prerm ui-slim] disabling display setup service"
            systemctl disable crankshaft-ui-slim-display-setup.service >/dev/null 2>&1 || true
        fi

        # If systemd handled stop, avoid further manual killing to prevent accidental TERM
        if ! command -v systemctl >/dev/null 2>&1 || [ ! -d /run/systemd/system ]; then
            if command -v start-stop-daemon >/dev/null 2>&1; then
                echo "[prerm ui-slim] stopping by name (no systemd)"
                start-stop-daemon --stop --quiet --oknodo --retry=TERM/5 --name crankshaft-ui-slim || true
                echo "[prerm ui-slim] stopping by exec (no systemd)"
                start-stop-daemon --stop --quiet --oknodo --retry=TERM/5 --exec /usr/bin/crankshaft-ui-slim || true
            fi
        fi

        echo "Crankshaft Slim UI pre-removal cleanup done"
    ;;
    upgrade)
        echo "[prerm ui-slim] phase=upgrade"
        # During upgrade, only stop the service gracefully
        if command -v systemctl >/dev/null 2>&1 && [ -d /run/systemd/system ]; then
            echo "[prerm ui-slim] stopping service"
            systemctl stop crankshaft-ui-slim.service >/dev/null 2>&1 || true
            echo "[prerm ui-slim] stopping display setup service"
            systemctl stop crankshaft-ui-slim-display-setup.service >/dev/null 2>&1 || true
        fi

        # Wait for process to exit (with timeout)
        for i in 1 2 3 4 5; do
            if ! pgrep -f '[c]rankshaft-ui-slim' >/dev/null 2>&1; then
                break
            fi
            echo "[prerm ui-slim] waiting for crankshaft-ui-slim to exit (attempt $i)"
            sleep 1
        done

        # Force kill only when systemd is absent (avoid double-stop TERM)
        if ! command -v systemctl >/dev/null 2>&1 || [ ! -d /run/systemd/system ]; then
            if command -v start-stop-daemon >/dev/null 2>&1; then
                echo "[prerm ui-slim] force kill by name (no systemd)"
                start-stop-daemon --stop --quiet --oknodo --retry=KILL/5 --name crankshaft-ui-slim || true
                echo "[prerm ui-slim] force kill by exec (no systemd)"
                start-stop-daemon --stop --quiet --oknodo --retry=KILL/5 --exec /usr/bin/crankshaft-ui-slim || true
            fi
        fi

        echo "Crankshaft Slim UI stopped for upgrade"
    ;;
esac

exit 0
